# scripts/localgo.service
[Unit]
Description=LocalGo - LocalSend v2.1 Protocol Implementation
Documentation=https://github.com/bethropolis/localgo
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=localgo
Group=localgo
ExecStart=/usr/local/bin/localgo-cli serve --quiet
Restart=always
RestartSec=5
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=localgo

# Environment configuration
EnvironmentFile=-/etc/localgo/localgo.env
WorkingDirectory=/var/lib/localgo

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/localgo /var/log/localgo
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateDevices=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryAccounting=true
MemoryMax=512M

[Install]
WantedBy=multi-user.target
