# scripts/localgo.service
# System service for localgo — installed to /usr/lib/systemd/system/localgo.service
#
# For package-manager installs (deb, rpm).
# Requires a 'localgo' system user and group (created by the postinstall script).
#
# Enable:  sudo systemctl enable --now localgo
# Logs:    sudo journalctl -u localgo -f
# Config:  /etc/localgo/localgo.env

[Unit]
Description=LocalGo - LocalSend v2.1 Protocol Implementation
Documentation=https://github.com/bethropolis/localgo
# Wait for the network stack to be up; multicast will fail without it.
After=network.target network-online.target
Wants=network-online.target
# Stop rebooting on repeated crashes: allow 3 starts per 60 s, then give up.
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=localgo
Group=localgo

ExecStart=/usr/bin/localgo serve --quiet --auto-accept


Restart=on-failure
RestartSec=5

# First start generates a TLS certificate — give it enough time.
# 10 s is too tight on ARM or slow I/O. 30 s is safe.
TimeoutStartSec=30
TimeoutStopSec=15

# Log everything to journald (stdout + stderr).
# null would silently swallow any stdout diagnostic output.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=localgo

# -----------------------------------------------------------------------
# Directory management
# -----------------------------------------------------------------------
StateDirectory=localgo
StateDirectoryMode=0750
ConfigurationDirectory=localgo
ConfigurationDirectoryMode=0750
WorkingDirectory=/var/lib/localgo

# -----------------------------------------------------------------------
# Environment
# -----------------------------------------------------------------------
Environment=LOCALSEND_DOWNLOAD_DIR=/var/lib/localgo/downloads
EnvironmentFile=-/etc/localgo/localgo.env

# -----------------------------------------------------------------------
# Security hardening
# -----------------------------------------------------------------------
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/localgo
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateDevices=true
CapabilityBoundingSet=

# Resource limits
LimitNOFILE=4096
LimitNPROC=256
MemoryAccounting=true
MemoryMax=256M

[Install]
WantedBy=multi-user.target
